# roles/haproxy_tls_certificate/templates/haproxy-certs-renew.service.j2
[Unit]
Description=Renew TLS certificates and prepare HAProxy combined PEM files
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile={{ haproxy_tls_certificate_haproxy_certs_dir }}/.env
ExecStart=/usr/local/bin/prepare-haproxy-certs.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=haproxy-certs
ProtectSystem=strict
ReadWritePaths={{ haproxy_tls_certificate_haproxy_certs_dir }} {{ haproxy_tls_certificate_lego_cert_dir }}
PrivateTmp=true
NoNewPrivileges=true
