# roles/k8s_oidc_github_pages/templates/k8s-oidc-github-pages.timer.j2
[Unit]
Description=https://github.com/meysam81/ansible-collections

[Timer]
OnCalendar={{ k8s_oidc_github_pages_timer_calendar }}

[Install]
WantedBy=multi-user.target
