# roles/haproxy_tls_certificate/files/haproxy-certs-renew.timer
[Unit]
Description=Daily TLS certificate renewal and HAProxy cert preparation

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
