# roles/cloudflare_ips/files/cloudflare-ips-update.service
[Unit]
Description=Update Cloudflare IP ranges for HAProxy
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-cloudflare-ips.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflare-ips
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true
