# roles/vmagent/templates/vmagent.service.j2
# {{ ansible_managed }}
[Unit]
Description=Victoria Metrics - VMAgent
Documentation=https://docs.victoriametrics.com/vmagent/
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Type=simple

ExecStartPre={{ vmagent_bin }} \
  -promscrape.config=/etc/vmagent/vmagent.yml \
  -promscrape.config.dryRun
ExecStart={{ vmagent_bin }} \
  -promscrape.config=/etc/vmagent/vmagent.yml \
  -remoteWrite.url={{ vmagent_remote_write_url }} \
{% if vmagent_remote_write_token_file %}
  -remoteWrite.bearerTokenFile={{ vmagent_remote_write_token_file }} \
{% endif %}
{% if vmagent_remote_write_username and vmagent_remote_write_password %}
  -remoteWrite.basicAuth.username={{ vmagent_remote_write_username }} \
  -remoteWrite.basicAuth.passwordFile=/etc/vmagent/remote-url-password \
{% endif %}
  -remoteWrite.tmpDataPath=/var/lib/vmagent/remote-write-data \
{% if vmagent_remote_write_force_vm_proto %}
  -remoteWrite.forceVMProto \
{% endif %}
{% if vmagent_instance_label %}
  -remoteWrite.label=instance={{ vmagent_instance_label }} \
{% endif %}
{% for key, value in vmagent_extra_labels.items() %}
  -remoteWrite.label={{ key }}={{ value }} \
{% endfor %}
{% for cli_flag in vmagent_extra_cli_flags %}
  {{ cli_flag }} \
{% endfor %}
  -httpListenAddr={{ vmagent_http_listen_addr }} \
  -influxListenAddr=""
ExecReload=/bin/kill -SIGHUP $MAINPID

# ── Identity ──────────────────────────────────────────────────────────
User=vmagent
Group=vmagent

# ── Capabilities ──────────────────────────────────────────────────────
CapabilityBoundingSet=
NoNewPrivileges=yes

# ── Filesystem isolation ──────────────────────────────────────────────
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ReadWritePaths=/var/lib/vmagent

# ── Kernel isolation ─────────────────────────────────────────────────
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes
ProtectProc=invisible
ProcSubset=pid

# ── Network isolation ────────────────────────────────────────────────
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK

# ── Syscall filtering ───────────────────────────────────────────────
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# ── Logging ─────────────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vmagent

# ── Misc hardening ──────────────────────────────────────────────────
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes
UMask=0077

# ── Resource limits ──────────────────────────────────────────────────
LimitNOFILE=65535
LimitNPROC=4096

# ── Restart policy ───────────────────────────────────────────────────
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
