# roles/tls_certificate/files/haproxy-certs-renew.timer
[Unit]
Description=Weekly TLS certificate renewal and HAProxy cert preparation

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
