# roles/vm_identity_gcp/files/vm-identity-gcp.timer
[Unit]
Description=Timer for VM Identity Token Retrieval

[Timer]
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=timers.target
